# dist/nhweb/DEBIAN/postinst
systemctl enable nhweb
service nhweb start
